# server/docker/pushr.sv.conf
[program:pushr]
command=/usr/local/sbin/pushr.run
redirect_stderr=true